# crontab
15 7 * * * osascript /path/to/morning.script